# workers-local/com.mgsr.releases-refresh.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.mgsr.releases-refresh</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/zsh</string>
        <string>-c</string>
        <string>source ~/.nvm/nvm.sh &amp;&amp; nvm use 22 &amp;&amp; node /Users/lior.dahan/AndroidStudioProjects/MGSRTeam/workers-local/releases-refresh.js</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/lior.dahan/AndroidStudioProjects/MGSRTeam/workers-local</string>

    <!-- Run daily at 10:00 local system time (set macOS timezone to Israel for exact 10:00 Israel time) -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>10</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/lior.dahan/AndroidStudioProjects/MGSRTeam/workers-local/logs/releases-stdout.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/lior.dahan/AndroidStudioProjects/MGSRTeam/workers-local/logs/releases-stderr.log</string>

    <key>RunAtLoad</key>
    <true/>
</dict>
</plist>
